# phone_diary/audiodiary.service
[Unit]
Description=Audiodiary

[Service]
ExecStart=bash ~/Documents/scripts/phone_diary/watchscript

[Install]
WantedBy=default.target
